# backend/scripts/yt-dlp-upgrade.timer
[Unit]
Description=Daily yt-dlp upgrade

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
